# deploy/systemd/life-engine-core.service
[Unit]
Description=Life Engine Core
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=life-engine
Group=life-engine
ExecStart=/usr/local/bin/life-engine-core --config /etc/life-engine/config.toml
Restart=on-failure
RestartSec=5

# Environment
Environment=LIFE_ENGINE_STORAGE_PASSPHRASE=
WorkingDirectory=/var/lib/life-engine

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/life-engine
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
